# backend/collect-pipeline.service
[Unit]
Description=App Store & Play Store Collection Pipeline Daemon
Documentation=file:///home/ubuntu/newapps/backend/README.md
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/newapps/backend

# 환경 변수 파일 (설치 후 주석 해제 필요)
# EnvironmentFile=/home/ubuntu/newapps/backend/.env

# 또는 개별 환경 변수 설정 (EnvironmentFile 대신 사용 가능)
# Environment="APP_DETAILS_DB_HOST=localhost"
# Environment="APP_DETAILS_DB_PORT=5432"
# Environment="APP_DETAILS_DB_NAME=app_details"
# Environment="APP_DETAILS_DB_USER=app_details"
# Environment="APP_DETAILS_DB_PASSWORD=your_password_here"

# 메인 실행 명령어
ExecStart=/home/ubuntu/newapps/backend/.venv/bin/python /home/ubuntu/newapps/backend/collect_full_pipeline.py --daemon

# Graceful Shutdown 설정
# - SIGTERM 전송 후 최대 5분 대기 (현재 작업 완료 대기)
# - 5분 후에도 종료되지 않으면 SIGKILL
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=300

# 자동 재시작 설정
# - on-failure: 비정상 종료 시에만 재시작 (exit 0은 재시작 안 함)
# - RestartSec: 재시작 전 30초 대기
Restart=on-failure
RestartSec=30

# 리소스 제한 (선택사항)
# MemoryMax=2G
# CPUQuota=80%

# 로그 설정
StandardOutput=journal
StandardError=journal
SyslogIdentifier=collect-pipeline

[Install]
WantedBy=multi-user.target
